# deploy/hummingbot.service
[Unit]
Description=Hummingbot factor MM (BTC perp testnet)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botuser
WorkingDirectory=/home/botuser/hummingbot
EnvironmentFile=/home/botuser/.hb-password.env
ExecStart=/home/botuser/factor-mm/deploy/run-hummingbot.sh --v2 factor_mm.yml --headless
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
